# freenas/atf/rc.local
#!/bin/sh
TERM="cons25"
PATH="${PATH}:/usr/local/bin:/usr/local/sbin"

clear

# Set the IP address we will use for testing
echo "Changing IP address for REST"
/sbin/ifconfig em0 192.168.0.15 255.255.255.0

sleep 5

if [ -e "/etc/netcli" ] ; then
  # Set the default password for 9.x
  echo "Setting ROOT password for REST"
  /etc/netcli reset_root_pw "testing"
else
  echo "Missing /etc/netcli...! RootPW not set!"
fi

sleep 5
