# freenas/atf/rc.local
#!/bin/sh
TERM="cons25"
PATH="${PATH}:/usr/local/bin:/usr/local/sbin"

clear

# Set the IP address we will use for testing
/sbin/ifconfig em0 192.168.0.15 255.255.255.0

if [ -e "/etc/netcli" ] ; then
  # Set the default password for 9.x
  /etc/netcli reset_root_pw "testing"
else
  echo "Missing /etc/netcli...! RootPW not set!"
fi

# Wait 5 minutes, then shutdown
sleep 300
shutdown -p now
